# etc/cron.hourly/dpkg-list
#!/bin/bash

set -e

t="$(tempfile)"

dpkg --get-selections > "$t"

cp "$t" "/etc/dpkg/packages.txt"

rm "$t"

if type -p etckeeper > /dev/null
then
    etckeeper commit -m 'dpkg list updated'
fi
